# python_libs_common/airootfs/etc/profile.d/set_pythonpath.sh
#!/bin/sh

export PYTHONPATH=/usr/lib/python-escape-common